# examples/remote-claude.service.example
# Example systemd user service the Linux bootstrap installs to
#   ~/.config/systemd/user/remote-claude.service
# It starts the tunnel and restarts it 15s after any exit.
#
# Manage it with:
#   systemctl --user enable --now remote-claude.service
#   systemctl --user disable --now remote-claude.service
#   journalctl --user -u remote-claude -f
#
# By default user services only run while you are logged in; to keep the
# tunnel up without an active session:
#   sudo loginctl enable-linger $USER

[Unit]
Description=remote-claude reverse SSH tunnel
After=network-online.target

[Service]
ExecStart=/usr/bin/ssh -N -o ExitOnForwardFailure=yes remote-claude
Restart=always
RestartSec=15

[Install]
WantedBy=default.target
